# etc/init.d/usbmuxd
#!/sbin/openrc-run

description="Socket daemon for the usbmux protocol used by Apple devices"
command="/usr/sbin/usbmuxd"
command_args="-U usbmux -f -z"
supervisor=supervise-daemon

depend() {
	need udev
}

start_pre() {
	if ! getent passwd usbmux >/dev/null; then
		eerror "User 'usbmux' does not exist. Please create it first."
		return 1
	fi
}
