# debian/postinst
#! /bin/bash

#   Add usbmuxd service files to runlevel/default

rc-update add usbmuxd-openrc default

#   Update services list

rc-update -u
